# com.user.ime-switcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ime-switcher</string>

    <key>ProgramArguments</key>
    <array>
        <!-- 下面这行改成你实际的可执行文件绝对路径 -->
        <string>/usr/local/bin/ime-switcher</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ime-switcher.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ime-switcher.err</string>
</dict>
</plist>
